# systemd/codexbar-smoke-tests.service
[Unit]
Description=CodexBar scheduled provider smoke tests
After=openclaw-gateway.service network-online.target
Wants=openclaw-gateway.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 %h/.local/bin/codexbar-linux.py scheduled-smoke-tests
